# movie-server/rootfs/etc/services.d/movie-server/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

MAIN_URL="$(bashio::config 'main_url')"
MAX_PAGES="$(bashio::config 'max_pages')"
TMDB_API_KEY="$(bashio::config 'tmdb_api_key')"
DOWNLOAD_DIR="$(bashio::config 'download_dir')"

INITIAL_PAGES="$(bashio::config 'initial_pages')"
if [ -z "${INITIAL_PAGES}" ]; then INITIAL_PAGES=2; fi

HD_KEYWORDS="$(bashio::config 'hd_keywords')"
HD_KEYWORDS="${HD_KEYWORDS#"${HD_KEYWORDS%%[![:space:]]*}"}"
HD_KEYWORDS="${HD_KEYWORDS%"${HD_KEYWORDS##*[![:space:]]}"}"
if [ -z "${HD_KEYWORDS}" ]; then
  HD_KEYWORDS="720p,1080p,HD,HDRip,WEB-DL,BluRay,Blu-Ray"
fi

K4_KEYWORDS="$(bashio::config 'k4_keywords')"
K4_KEYWORDS="${K4_KEYWORDS#"${K4_KEYWORDS%%[![:space:]]*}"}"
K4_KEYWORDS="${K4_KEYWORDS%"${K4_KEYWORDS##*[![:space:]]}"}"
if [ -z "${K4_KEYWORDS}" ]; then
  K4_KEYWORDS="2160p,4k,4K,UHD"
fi

EMBY_URL="$(bashio::config 'emby_url')"
EMBY_URL="${EMBY_URL#"${EMBY_URL%%[![:space:]]*}"}"
EMBY_URL="${EMBY_URL%"${EMBY_URL##*[![:space:]]}"}"

EMBY_API_KEY="$(bashio::config 'emby_api_key')"

EMBY_PATH_PREFIX="$(bashio::config 'emby_path_prefix')"
EMBY_PATH_PREFIX="${EMBY_PATH_PREFIX#"${EMBY_PATH_PREFIX%%[![:space:]]*}"}"
EMBY_PATH_PREFIX="${EMBY_PATH_PREFIX%"${EMBY_PATH_PREFIX##*[![:space:]]}"}"

cat > /app/.env <<EOF
MAIN_URL=${MAIN_URL}
PORT=3001
MAX_PAGES=${MAX_PAGES}
INITIAL_PAGES=${INITIAL_PAGES}
TMDB_API_KEY=${TMDB_API_KEY}
DOWNLOAD_DIR=${DOWNLOAD_DIR}
HD_KEYWORDS=${HD_KEYWORDS}
K4_KEYWORDS=${K4_KEYWORDS}
EMBY_URL=${EMBY_URL}
EMBY_API_KEY=${EMBY_API_KEY}
EMBY_PATH_PREFIX=${EMBY_PATH_PREFIX}
EOF

cd /app || exit 1
exec node main.js
