# movie-server/rootfs/etc/services.d/movie-server/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

cat > /app/.env <<EOF
MAIN_URL=$(bashio::config 'main_url')
PORT=3001
MAX_PAGES=$(bashio::config 'max_pages')
INITIAL_PAGES=$(bashio::config 'initial_pages')
TMDB_API_KEY=$(bashio::config 'tmdb_api_key')
DOWNLOAD_DIR=$(bashio::config 'download_dir')
HD_KEYWORDS=$(bashio::config 'hd_keywords')
K4_KEYWORDS=$(bashio::config 'k4_keywords')
EOF

cd /app || exit 1
exec node main.js
